# build-pack/_CPack_Packages/Linux/TGZ/aegisbpf-0.1.0-Linux/usr/lib/systemd/system/aegisbpf.service
[Unit]
Description=AegisBPF Runtime Security Agent
After=network.target
ConditionPathExists=/sys/fs/bpf
ConditionPathExists=/sys/fs/cgroup/cgroup.controllers

[Service]
Type=simple
EnvironmentFile=-/etc/default/aegisbpf
ExecStartPre=/bin/sh -c '[ -n "${AEGIS_POLICY:-}" ] && [ -f "${AEGIS_POLICY}" ] && /usr/local/bin/aegisbpf policy apply "${AEGIS_POLICY}" --reset'
ExecStart=/usr/local/bin/aegisbpf run ${AEGIS_MODE} ${AEGIS_LOG}
Restart=on-failure
RestartSec=2
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
